# sleeper.service
[Unit]
Description=Simple sleeper Logger

[Service]
ExecStart=%h/bin/sleeper 60
Restart=no

[Install]
WantedBy=multi-user.target