# sleeper.service
[Unit]
Description=Simple sleeper Logger

[Service]
ExecStart=%h/bin/sleeper
Restart=on-failure

[Install]
WantedBy=multi-user.target
